# scripts/nas-monitor.service
[Unit]
Description=NAS Monitor Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/nas_monitor
Environment="PATH=/opt/nas_monitor/.venv/bin:/root/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin"
ExecStart=/opt/nas_monitor/run.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
